# temp_services/phase3.service
[Unit]
Description=Configuración Proxmox Workstation - Fase 3 (Escritorios)
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
WorkingDirectory=/home/zenosama
ExecStart=/opt/pve-setup/setup-pve-workstation-phase3.sh

[Install]
WantedBy=multi-user.target