# custom_services/autologin.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin zenosama --noclear tty1
