# custom_services/autologin.service
[Service]
Type=oneshot
ExecStart=
ExecStart=-/sbin/agetty --autologin zenosama --noclear %I $TERM
