# custom_services/phase2.service
[Unit]
Description=Configuración Proxmox Workstation - Fase 2 (instalación de paquetes principales)
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
WorkingDirectory=/home/zenosama
ExecStart=/opt/pve-setup/setup-pve-workstation-phase2.sh

[Install]
WantedBy=multi-user.target